# scripts/deploy/intranet.socket
[Unit]
Description=gunicorn socket for intranet

[Socket]
ListenStream=/run/intranet.sock

[Install]
WantedBy=sockets.target
